# workshop/supervisor/examiner.conf
[program:examiner]
process_name=examiner
command=/home/eduk8s/examiner/start-examiner.sh
stdout_logfile=/proc/1/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
